# infra/supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/stdout
logfile_maxbytes=0

[program:express]
command=node --import tsx backend/src/api_server.ts
directory=/app
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=NODE_ENV="production",PORT="8080"

; FastAPI is currently disabled. The ADK workflow runs in-process via Node.js.
; autostart MUST be a strict boolean (true/false).
; Do NOT use %(ENV_ADK_ENABLED)s here — supervisord rejects non-boolean values like "controlled".
[program:fastapi]
command=python3 -m uvicorn backend.python.main:app --host 0.0.0.0 --port 8000
directory=/app
autostart=false
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=PYTHONPATH="/app"

